# dockermultiarch.service
[Unit]
Description=DockerMultiArchSupportBuild

[Service]
ExecStart=/usr/bin/dockermultiarch

[Install]
WantedBy=multi-user.target 